# docker/elk/conf/run.conf
[program:run]
command=/root/app/run.sh
priority=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/run.log
stdout_logfile_maxbytes = 10MB
