# docker/elk/conf/run.conf
[program:jdbc-elk]
command=/root/app/run.sh
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/run-out.log
stderr_logfile=/var/log/supervisor/run-err.log
